# chatbot.service
[Unit]
Description=Chatbot API Django Application
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/chatbot_api
Environment="PATH=/var/www/chatbot_api/venv/bin"
Environment="DJANGO_SETTINGS_MODULE=chatbot_api.production_settings"
Environment="DJANGO_SECRET_KEY=your-secret-key-here"
ExecStart=/var/www/chatbot_api/venv/bin/gunicorn --workers 3 --bind 0.0.0.0:8000 --log-level debug --access-logfile /var/log/gunicorn/access.log --error-logfile /var/log/gunicorn/error.log chatbot_api.wsgi:application
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target 